# config/com.user.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- Label for the job -->
  <key>Label</key>
  <string>com.user.backup</string>

  <!-- Program to run -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/path/to/backup-system/scripts/backup.sh</string>
  </array>

  <!-- Schedule: Run every day at 2am and 2pm -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>2</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>14</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>

  <!-- WatchPaths: Monitor changes to the config files -->
  <key>WatchPaths</key>
  <array>
    <string>/path/to/backup-system/config/restic.env</string>
  </array>

  <!-- Log paths for standard output and error -->
  <key>StandardOutPath</key>
  <string>/Users/youruser/Library/Logs/backup.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/youruser/Library/Logs/backup.err</string>

  <!-- Automatically run the job when loaded -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Run when system wakes from sleep -->
  <key>StartOnMount</key>
  <true/>

  <!-- Prevent too frequent runs if system wakes up multiple times -->
  <key>ThrottleInterval</key>
  <integer>3600</integer>

  <!-- Ensure job runs even if system was asleep at scheduled time -->
  <key>StartInterval</key>
  <integer>86400</integer>
</dict>
</plist>
